# supervisor/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch -Des.config=/data/elasticsearch.yml
redirect_stderr=true
autorestart=true
