# supervisor/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
redirect_stderr=true
autorestart=true
